# autodev.service
[Unit]
Description=Human-AI Autonomous Development Daemon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/human-ai
ExecStart=/bin/bash /home/ubuntu/autodev.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
