# systemd-services/liveea.service
[Unit]
Description=Human-AI EA Trading Agent
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/yahwehatwork/human-ai
ExecStart=/home/yahwehatwork/human-ai/.venv/bin/python3 /home/yahwehatwork/human-ai/liveea.py
Restart=always
RestartSec=10
StandardOutput=append:/home/yahwehatwork/human-ai/data/logs/liveea.service.log
StandardError=append:/home/yahwehatwork/human-ai/data/logs/liveea.service.err

[Install]
WantedBy=multi-user.target
